# scripts/supervisor/alice.conf
[program:relay-alice]
command="./root/scripts/alice.sh"
stderr_logfile=/var/log/bitgreen-relay-alice-error.log
stdout_logfile=/var/log/bitgreen-relay-alice.log
