# scripts/supervisor/bob.conf
[program:relay-bob]
command="./root/polkadot/scripts/alice.sh"
stderr_logfile=/var/log/bitgreen-relay-bob-error.log
stdout_logfile=/var/log/bitgreen-relay-bob.log
