# scripts/supervisor/bob.conf
[program:relay-bob]
command=~/polkadot/polkadot --chain ~/rococo-local-cfde.json --bob --tmp
stderr_logfile = /var/log/bitgreen/relay-bob-error.log
stdout_logfile = /var/log/bitgreen/relay-bob.log
