# scripts/supervisor/parachain.conf
[program:parachain]
command="/root/polkadot/scripts/parachain.sh"
stderr_logfile=/var/log/bitgreen-parachain-error.log
stdout_logfile=/var/log/bitgreen-parachain.log
