# scripts/supervisor/parachain.conf
[program:parachain]
command="/root/bitg-parachain/target/release/bitg-parachain --dev --alice --validator -- --chain ~/rococo-local-cfde.json"
stderr_logfile=/var/log/bitgreen/parachain-error.log
stdout_logfile=/var/log/bitgreen/parachain.log
